# files/voltdb.conf
## Upstart config file (use 'start voltdb', 'stop voltdb')
author "Sean Van Osselaer"
description "Start VoltDB"

start on (local-filesystems and net-device-up IFACE!=lo)

exec voltdb create
